# .config/systemd/user/dotfiles-sync.service
[Unit]
Description=Sync dotfiles and apply stow symlinks
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c '[[ -x "$HOME/.config/dotfiles/dotfiles-sync.sh" ]] && "$HOME/.config/dotfiles/dotfiles-sync.sh"'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
