# .config/systemd/user/vscodium-sync.service
[Unit]
Description=Sync VSCodium extensions from extensions.list
After=After=dotfiles-sync.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c '[[ -x "$HOME/.config/VSCodium/install-extensions.sh" ]] && "$HOME/.config/VSCodium/install-extensions.sh"'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
